# salt/pwm-server/pwm-server.conf
description "pwm-server job"

start on startup
stop on runlevel [016]

respawn

exec /srv/pwm-server/venv/bin/pwm-server --config /srv/pwm-server/config.py
